# systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- manuel' --noclear --skip-login - $TERM
